# src/main/udev/rules.d/makerbot.rules
KERNEL=="tty*", SUBSYSTEM=="tty", RUN+="conveyor-udev"
