# data/conveyor.conf
# conveyor - daemon for controlling 3d printers

description	"3D printing control system"
author "MakerBot dev team <dev@makerbot.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec sudo -u conveyor /usr/bin/conveyord --nofork -c /etc/conveyor.conf
