# data/conveyor.conf
# conveyor - daemon for controlling 3d printers

description	"3D printing control system"
author "MakerBot dev team <dev@makerbot.com>"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn
respawn limit 10 20


exec start-stop-daemon -- start -c conveyor --exec python /usr/share/conveyor/conveyor_service -c /etc/conveyor.conf
