# data/conveyor.conf
# conveyor - daemon for controlling 3d printers

description	"3D printing control system"
author "MakerBot dev team <dev@makerbot.com>"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn
respawn limit 10 20

#pre-start script
# consider picking up a /etc/default/conveyor for location of conf file here
#end script

exec /usr/bin/conveyord -c /etc/conveyor.conf
